# sp-analyzer/src/sp-analyzer/etc/cron.d/sp-analyzer.cron
# storage-processor analyzer
*/3 * * * * root /usr/local/bin/sp-analyzer.py >/dev/null 2>&1
